# infra/systemd/animal-welfare-weekly-crawl.timer
[Unit]
Description=Run weekly browser crawl every Sunday at 03:00

[Timer]
OnCalendar=Sun *-*-* 03:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
